# voice-agent/aura_voice.service
[Unit]
Description=AURA Voice Agent — AI-powered apartment assistant
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=root
WorkingDirectory=/config/aura/voice-agent
ExecStartPre=/bin/sh -c 'test -f /config/aura/.env && test -d /config/aura/.venv'
ExecStart=/config/aura/.venv/bin/python aura_voice.py
Restart=on-failure
RestartSec=10
TimeoutStopSec=15
# WatchdogSec=120  # Enable after adding sd_notify keepalive to aura_voice.py

# Environment
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=/config/aura/.env

# Resource limits
MemoryMax=512M
MemoryHigh=384M
CPUQuota=80%

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/config/aura
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=aura-voice

[Install]
WantedBy=multi-user.target
