# clap-trigger/clap_service.service
[Unit]
Description=AURA Clap Detection Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/pi/aura/clap-trigger
ExecStart=/home/pi/aura/.venv/bin/python clap_listener.py
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
